# files/etc/service/sshd/run
#!/bin/sh
set -e

if [ ! -d /var/run/sshd ]; then
   mkdir /var/run/sshd
   chmod 0755 /var/run/sshd
fi

exec /usr/sbin/sshd -D
