# rootfs/etc/service/incrond/run
#!/bin/sh
set -e

exec /usr/sbin/incrond -f /etc/incron.conf -n
